# cfg/etc/periodic/15min/kubelet_checkin
if [ ! -f /var/run/kubelet.pid ]; then
	exit 0
fi

if [ -f /var/run/kubelet.exitcode ]; then
	service kubelet restart
else
	rm -fv /etc/periodic/15min/kubelet_checkin
fi