# cfg/etc/periodic/15min/persist_var_lib_kubelet
rsync -rv --no-links /var/lib/kubelet/pki /mnt/persist/$HOSTNAME/var/lib/kubelet/ 