# cfg/etc/periodic/15min/persist_var_lib
rsync -rav --no-links --exclude '*snapshots*' --exclude '*.sock' /var/lib/ /mnt/persist/$HOSTNAME/var/lib/