# cfg/etc/periodic/15min/persist_var_lib
if [ -f /etc/periodic/15min/kubelet_checkin ]; then
	exit 0
fi

rsync -ravh  --ignore-errors --delete --no-links --exclude '*snapshots*' --exclude '*.sock' /var/lib/ /mnt/persist/$HOSTNAME/var/lib/